# infra/systemd/fitness-token-refresh.service
[Unit]
Description=Refresh Strava/Withings OAuth tokens for OpenClaw fitness skill
After=openclaw-gateway.service

[Service]
Type=oneshot
ExecStart=/usr/bin/sg docker -c /home/ophirsw/Documents/claude/nuc-vector-orchestrator/scripts/refresh-fitness-tokens.sh
TimeoutSec=60
