# infra/systemd/nuc-agent-loop.service
[Unit]
Description=NUC Agent Loop — PM, Coder, Tester, CISO dispatch for Project Vector
After=network-online.target

[Service]
Type=notify
NotifyAccess=main
ExecStart=/usr/bin/python3 -m apps.control_plane.agent_loop
Restart=always
RestartSec=30
WatchdogSec=2400
Environment=HOME=/home/ophirsw
Environment=PATH=/home/ophirsw/.nvm/versions/node/v24.14.0/bin:/home/ophirsw/.local/bin:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=/home/ophirsw/Documents/claude/nuc-vector-orchestrator

[Install]
WantedBy=default.target
