# infra/systemd/openclaw-backup.timer
[Unit]
Description=OpenClaw hourly backup timer

[Timer]
OnCalendar=hourly
RandomizedDelaySec=3min
Persistent=true

[Install]
WantedBy=timers.target
