# infra/systemd/vector-bridge.service
[Unit]
Description=Vector HTTP-to-gRPC Bridge Server
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=%i
WorkingDirectory=%h/Documents/claude/nuc-vector-orchestrator
ExecStart=/usr/bin/python3 -m apps.vector.bridge.server --port 8081 --host 127.0.0.1
Restart=on-failure
RestartSec=10
Environment=PYTHONPATH=%h/Documents/claude/nuc-vector-orchestrator
Environment=VECTOR_SERIAL=0dd1cdcf
Environment=VECTOR_BRIDGE_PORT=8081
Environment=VECTOR_BRIDGE_HOST=127.0.0.1
# Also binds to 172.17.0.1 (Docker bridge) so OpenClaw container can reach the bridge

# Security hardening
NoNewPrivileges=yes
ProtectSystem=strict
ReadWritePaths=%h/Documents/claude/nuc-vector-orchestrator
PrivateTmp=yes

[Install]
WantedBy=default.target
