# infra/vector/vector-streamer.service
# vector-streamer.service — Mic/engine proxy + TCP bridge to NUC
#
# Uses systemd's native ordering to insert socket proxies between services:
#
#   vic-anim → vector-streamer → vic-engine → vic-cloud
#
# Before= ensures systemd won't start vic-engine or vic-cloud until
# vector-streamer has finished setting up its socket proxies (Type=forking,
# so "ready" = startup script exits after verifying proxies are live).
#
# If vector-streamer is not enabled, Before= has no effect — services
# boot normally without proxies.
#
# Install:
#   ssh vector 'mount -o remount,rw /'
#   scp vector-streamer-start.sh root@vector:/anki/bin/
#   scp vector-streamer.service root@vector:/etc/systemd/system/
#   ssh vector 'chmod +x /anki/bin/vector-streamer-start.sh'
#   ssh vector 'systemctl daemon-reload && systemctl enable vector-streamer'
#   ssh vector 'mount -o remount,ro /'
#
# Logs:
#   ssh vector 'cat /tmp/streamer.log'
#   ssh vector 'journalctl -u vector-streamer'

[Unit]
Description=Vector streamer — mic/engine proxy + TCP bridge to NUC

# Start after vic-anim (we need its sockets to exist before proxying them)
After=vic-anim.service
Wants=vic-anim.service

# Hold vic-engine and vic-cloud until our proxies are ready.
# This is the key line — systemd delays their startup until our
# forking script exits, meaning socket proxies are in place.
Before=vic-engine.service vic-cloud.service

[Service]
Type=forking
PIDFile=/tmp/vector-streamer.pid
ExecStart=/anki/bin/vector-streamer-start.sh
ExecStop=/anki/bin/vector-streamer-start.sh --stop

Restart=on-failure
RestartSec=10

# Startup script waits for vic-anim sockets + verifies proxy setup
TimeoutStartSec=90
TimeoutStopSec=30

[Install]
WantedBy=anki-robot.target
