# infra/openclaw/systemd/openclaw-backup.service
[Unit]
Description=Backup OpenClaw config to git
After=network.target

[Service]
Type=oneshot
WorkingDirectory=/home/ophirsw/Documents/claude/nuc-orchestrator
ExecStart=/usr/bin/bash /home/ophirsw/Documents/claude/nuc-orchestrator/scripts/backup-openclaw.sh
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-backup
