# infra/openclaw/systemd/openclaw-ade.service
[Unit]
Description=OpenClaw Gateway - Ade (Work Bot)
Documentation=https://openclaw.ai
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
WorkingDirectory=/home/ophirsw/Documents/claude/openclaw
ExecStart=/usr/bin/sg docker -c "/home/ophirsw/Documents/claude/openclaw/openclaw-ade-run.sh"
ExecStop=/usr/bin/sg docker -c "docker stop openclaw-ade"
Restart=always
RestartSec=10
TimeoutStartSec=90
TimeoutStopSec=30
StandardOutput=journal
StandardError=journal
SyslogIdentifier=openclaw-ade

[Install]
WantedBy=default.target
